# powershell-lolbins/systemd/powershell-lolbins-digest.timer
[Unit]
Description=Run the Wazuh "PowerShell / LOLBins" weekly summary (Mon 07:00)

[Timer]
OnCalendar=Mon *-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
