# dns-threat/systemd/dns-threat-digest.service
[Unit]
Description=Wazuh "DNS threats" weekly summary email
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Run as a dedicated, unprivileged service user (see README) -- NOT root.
# Create it once:  useradd --system --no-create-home --shell /usr/sbin/nologin wazuh-recipes
User=wazuh-recipes
Group=wazuh-recipes
RuntimeMaxSec=600
# Secrets (INDEX_PASSWORD, ...) and overrides live here, chmod 600, owned by the user.
EnvironmentFile=/etc/dns-threat-digest.env
# Edit --recipients to taste, or set RECIPIENTS in the env file and drop the flag.
ExecStart=/opt/dns-threat-digest/venv/bin/python \
          /opt/dns-threat-digest/dns-threat_digest.py \
          --recipients security@example.com
