# bruteforce-success/systemd/bruteforce-success-digest.service
[Unit]
Description=Wazuh "Brute-force success" weekly summary email
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Run as a dedicated, unprivileged service user (see README) -- NOT root.
# Create it once:  useradd --system --no-create-home --shell /usr/sbin/nologin wazuh-recipes
User=wazuh-recipes
Group=wazuh-recipes
RuntimeMaxSec=600
# Secrets (INDEX_PASSWORD, ...) and overrides live here, chmod 600, owned by the user.
EnvironmentFile=/etc/bruteforce-success-digest.env
# Edit --recipients to taste, or set RECIPIENTS in the env file and drop the flag.
ExecStart=/opt/bruteforce-success-digest/venv/bin/python \
          /opt/bruteforce-success-digest/bruteforce-success_digest.py \
          --recipients security@example.com
